# debian/prerm
#!/bin/sh
set -e
systemctl stop    picam-frontend 2>/dev/null || true
systemctl disable picam-frontend 2>/dev/null || true
